# ansible/roles/dPlanet/templates/etc/systemd/system/certbot-renewal.service
[Unit]
Description=Certbot Renewal

[Service]
ExecStart=/usr/bin/certbot renew --post-hook "docker restart $(docker ps -q --filter="name=nginx")"
